# nemoi-stt.service
[Unit]
  Description=Nemoi STT (for Telegram)

[Service]
  ExecStart=
  Type=simple

[Install]
  WantedBy=multi-user.target
